# ops/gcp/grt-live-scorer.service
[Unit]
Description=GRT Reliability Tracker live delay scorer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/grt-reliability-tracker
EnvironmentFile=/opt/grt-reliability-tracker/.env
# Scores every 5 minutes and uploads live-predictions.json to GCS_LIVE_BUCKET
# (set in .env) so the statically hosted dashboard can fetch it. Needs the
# trained model under data/analysis/models_live and parsed static GTFS under
# data/parsed_static_gtfs (see ops/gcp/README.md).
ExecStart=/opt/grt-reliability-tracker/collector/.venv/bin/python -u /opt/grt-reliability-tracker/analysis/predict_live.py --interval-seconds 300
Restart=always
RestartSec=30
User=grtcollector
Group=grtcollector

[Install]
WantedBy=multi-user.target
