# ops/gcp/grt-weather-forecast.service
[Unit]
Description=Collect ECCC weather forecast snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/grt-reliability-tracker
EnvironmentFile=/opt/grt-reliability-tracker/.env
ExecStart=/opt/grt-reliability-tracker/collector/.venv/bin/python -u /opt/grt-reliability-tracker/collector/collect_weather_forecasts.py
User=grtcollector
Group=grtcollector
