# ops/launchd/com.enkailiu.grt-reliability.daily-parse.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.enkailiu.grt-reliability.daily-parse</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/enkailiu/grt-reliability-tracker/collector/.venv/bin/python</string>
    <string>-u</string>
    <string>/Users/enkailiu/grt-reliability-tracker/collector/run_local_parse.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/enkailiu/grt-reliability-tracker</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/enkailiu</string>
    <key>PATH</key>
    <string>/Users/enkailiu/grt-reliability-tracker/collector/.venv/bin:/Users/enkailiu/google-cloud-sdk/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/enkailiu/grt-reliability-tracker/logs/daily-parse.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/enkailiu/grt-reliability-tracker/logs/daily-parse.err.log</string>
</dict>
</plist>
